# mysql-dumper/rc.local
#!/bin/sh
set -eu
if [ -n "${GOOGLE_APPLICATION_CREDENTIALS:-}" ]
then
    gcloud auth activate-service-account \
           --key-file=${GOOGLE_APPLICATION_CREDENTIALS}
fi
