# arachni-api/rc.local
#!/bin/sh

set -eu

/app/bin/postgres_setup
