# sensu/base/rc.local
#!/bin/sh
set -eu
exec /opt/sensu/sensu.env run-parts --verbose --exit-on-error /etc/rc.local.d
