# inspector/rc.local
#!/bin/sh
set -eu
gid=$(ls -ldn /var/run/docker.sock | cut -d' ' -f 4)
groupmod -g $gid docker
adduser --system --no-create-home app
getent group docker
getent passwd app
